# ansible/roles/k8s_firecracker/tasks/scripts/devmapper_reload.service
[Unit]
Description=Devmapper reload script

[Service]
ExecStart=/usr/local/bin/reload.sh

[Install]
WantedBy=multi-user.target